# foundrygate.service
[Unit]
Description=FoundryGate – Local AI Gateway for OpenClaw and other clients
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=foundrygate
Group=foundrygate
WorkingDirectory=/opt/foundrygate
EnvironmentFile=/opt/foundrygate/.env
Environment=FOUNDRYGATE_DB_PATH=/var/lib/foundrygate/foundrygate.db
ExecStart=/opt/foundrygate/venv/bin/uvicorn \
    foundrygate.main:app \
    --host 127.0.0.1 \
    --port 8090 \
    --workers 1 \
    --log-level info
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/foundrygate
PrivateTmp=true

[Install]
WantedBy=multi-user.target
