# docs/examples/foundrygate-auto-update.cron
# Example cron entry for controlled FoundryGate auto-update checks.
# Runs every day at 03:30 and only applies updates that the runtime marks as eligible.
30 3 * * * foundrygate /usr/local/bin/foundrygate-auto-update --apply >>/var/log/foundrygate-auto-update.log 2>&1
